# os/linux/etc/rc.local
#!/bin/sh

loadkeys /etc/my_keymaps_tty

exit 0
